# scripts/calendar-backend.service.tmpl
[Unit]
Description=Family Calendar backend (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@USER@
WorkingDirectory=@REPO@
ExecStart=@REPO@/backend/.venv/bin/uvicorn app.main:app --app-dir @REPO@/backend --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
